# Private.internal/Ascend_Infra/Ascend Launch Site.0/ascend.service
[Unit]
Description=Ascend AI Boot Service
After=network.target

[Service]
Type=simple
ExecStart=/home/you/Ascend/Neo
Restart=on-failure

[Install]
WantedBy=multi-user.target